# imageroot/systemd/user/grantora-retention.service
[Unit]
Description=Grantora audit and usage retention
After=grantora-api.service
Requires=grantora-api.service

[Service]
Type=oneshot
ExecStart=grantora-operations retention