# imageroot/systemd/user/grantora-user-sync.service
[Unit]
Description=Grantora LDAP user sync
After=grantora-api.service
Requires=grantora-api.service

[Service]
Type=oneshot
ExecStart=grantora-users sync